# crontab
# m h D M W command
  * * * * * kubectl version --client
